# deps/podman/conf/mypodmand.service
[Unit]
Description=Podman API Service
After=network.target

[Service]
Type=simple
User=root
# Environment=PATH=/opt/mypodman/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/mypodman/podman-server system service --time=0 unix:///tmp/mypodman.sock
ExecStop=/opt/mypodman/podman-server stop --all
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
