# deps/podman/conf/mypodmand.service
[Unit]
Description=Podman API Service
After=network.target

[Service]
Type=simple
User=root
Environment=PATH=/opt/mypodman:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/mypodman/podman system service --time=0 unix:///tmp/mypodman.sock
ExecStop=/usr/bin/podman stop --all
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
