# codespace/images/dev/rootfs/etc/s6/s6-rc.d/atuin-service/run
#!/opt/sb/profile/s6/bin/execlineb -P

# longrun run script (execline). Redirect output to the log, export the atuin
# server settings inline, then exec the atuin sync server.
# See `atuin server start --help` for the full list of tunables.
s6-envdir -Lf -- /run/s6/container_environment

fdmove -c 2 1
redirfd -w 1 /var/log/atuin-service.log

importas -i ATUIN_DB_URI ATUIN_DB_URI

# Hardcoded atuin server environment. Replace ATUIN_DB_URI with your real
# Supabase (or Postgres) connection string.
export ATUIN_HOST "127.0.0.1"
export ATUIN_PORT "2000"
export ATUIN_OPEN_REGISTRATION "false"
# export ATUIN_DB_URI "${ATUIN_DB_URI}"
exec /opt/sb/store/atuin/bin/atuin server start
