# dotfiles/s6/s6-rc.d/sshd/run
#!/opt/sb/profile/s6/bin/execlineb -P

# longrun run script (execline). Load the container environment dumped by the
# init, then exec sshd with the wanted port.
s6-envdir -Lf -- /run/s6/container_environment

fdmove -c 2 1
redirfd -w 1 /var/log/sshd.stdout.log

# Read SSHD_PORT from the now-loaded environment into an execline variable so it
# can be substituted into the exec line below (execline expands $vars itself).
importas -i SSHD_PORT SSHD_PORT

# Tighten host-key permissions and make sure the log dir exists.
foreground { mkdir -p /var/log }
elglob -0 keys /opt/devspace/dotfiles/sshd/host-key/*
foreground { chmod 600 $keys }

# Override the listen port on the command line (-o) instead of editing the
# config file, so the shipped sshd_config.conf stays untouched.
# https://github.com/un-def/openssh-static-build/blob/master/run-sshd.sh#L30
exec /opt/sb/store/openssh_gssapi/bin/sshd -D
  -o Port=${SSHD_PORT}
  -f /opt/devspace/dotfiles/sshd/sshd_config.conf
  -E /var/log/sshd.log
