# codespace/images/agent/rootfs/etc/s6/s6-rc.d/agent-service/run
#!/opt/sb/profile/s6/bin/execlineb -P

s6-envdir -Lf -- /run/s6/container_environment

fdmove -c 2 1

importas WORKSPACE_ROOT_HOST WORKSPACE_ROOT_HOST

exec /app/.venv/bin/python -m codespace.agent serve
  --host 0.0.0.0
  --port 8001
  --workspace-root-host ${WORKSPACE_ROOT_HOST}
  --podman-uri unix:///tmp/podmanxd.sock
